# deploy/allinone/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:postgres]
command=su postgres -c "%(ENV_PGBIN)s/postgres -D /data/postgres -c listen_addresses=127.0.0.1 -p 5432"
environment=PGBIN="/usr/lib/postgresql/15/bin"
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis]
command=redis-server --dir /data/redis --appendonly yes --bind 127.0.0.1
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:minio]
command=/usr/local/bin/minio server /data/minio --address %(ENV_MINIO_ADDRESS)s
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:init]
command=/usr/local/bin/init-app.sh
priority=20
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:api]
command=/opt/venv/bin/gunicorn apps.api.main:app -w %(ENV_API_WORKERS)s -k uvicorn.workers.UvicornWorker --bind 127.0.0.1:8000 --timeout 120 --graceful-timeout 30
directory=/workspace
priority=30
startretries=30
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
command=/opt/venv/bin/celery -A apps.api.tasks.celery_app worker -Q transcoding -c %(ENV_TRANSCODING_CONCURRENCY)s --loglevel=warning
directory=/workspace
priority=30
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:email_worker]
command=/opt/venv/bin/celery -A apps.api.tasks.celery_app worker -Q email_high,email_low -c %(ENV_EMAIL_CONCURRENCY)s --loglevel=warning
directory=/workspace
priority=30
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:beat]
command=/opt/venv/bin/celery -A apps.api.tasks.celery_app beat --loglevel=warning -s /tmp/celerybeat-schedule
directory=/workspace
priority=30
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:web]
command=node /app/web/server.js
directory=/app/web
environment=PORT="3000",HOSTNAME="0.0.0.0",NODE_ENV="production"
priority=30
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
priority=40
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
